# flask-app.service
[Unit]
Description=Flask app service
After=network.target

[Service]
WorkingDirectory=/home/ksmirnov/app
ExecStart=/home/ksmirnov/app/venv/bin/python -m gunicorn -b 0.0.0.0:80 -w 4
Restart=Always

[Install]
WantedBy=multi-user.target